# macos-maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.wind.weeklymaintenance</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-l</string>
      <string>~/Library/Scripts/macos-maintenance.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key>
      <integer>1</integer> <!-- 0 = Sunday -->
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </dict>
</plist>